# chfs.socket
[Unit]
Conflicts=chfs.service

[Socket]
ListenStream=80

[Install]
WantedBy=sockets.target
