# chfs/chfs.socket
[Unit]
Description="To Start Cut HTTP File Server"

[Socket]
ListenStream=80

[Install]
WantedBy=sockets.target
